# ttyd-proxy.service
[Unit]
Description=ttyd Telegram Mini App Auth Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/xithus/.hermes/ttyd-proxy
ExecStart=/home/linuxbrew/.linuxbrew/bin/node /home/xithus/.hermes/ttyd-proxy/server.js
EnvironmentFile=/home/xithus/.hermes/ttyd-proxy/.env
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
